# installation/deps/sspl-ll-test.service
[Unit]
Description=SSPL Low Level Process
Wants=openhpid.service drivemanager.service dcs-collector.service rabbitmq-server.service halond.service
After=rabbitmq-server.service dcs-collector.service halond.service syslog.target

[Service]
Type=notify
ExecStartPre=/usr/bin/mkdir /var/run/sspl_ll
ExecStart=
ExecStopPost=/usr/bin/rm -Rf /var/run/sspl_ll
PIDFile=/var/run/sspl_ll/sspl_ll_d.pid
User=root


[Install]
WantedBy=multi-user.target
