# installation/deps/sspl-ll.service.d/test-overrides.conf
[Unit]
Wants=rabbitmq-server.service
After=rabbitmq-server.service syslog.target

[Service]
ExecStart=/root/low-level/framework/sspl_ll_d -c /etc/sspl_ll.conf
User=root

